# rclone/wallpapers/driveBiSync_wallpapers.cron
@daily ~/dotfiles/rclone/wallpapers/driveBiSync_wallpapers.sh
