# rclone/wallpapers/driveBiSync_wallpapers.cron
@hourly anacron -s -t "$HOME/dotfiles/rclone/wallpapers/driveBiSync_wallpapers.anacron" -S "$HOME/.local/var/spool/anacron"
