# crontabs/driveBiSync.cron
@reboot /home/blue/dotfiles/rclone/driveBiSync.sh
*/10 * * * * /home/blue/dotfiles/rclone/driveBiSync.sh
