# rclone/obsidian_vaults/driveBiSync_obsidian-vaults.cron
#!/usr/bin/env bash
OBSIDIAN_CRON=$"@reboot $SCRIPT_DIR/obsidian_vaults/driveBiSync_obsidian-vaults.sh
*/15 * * * * $SCRIPT_DIR/obsidian_vaults/driveBiSync_obsidian-vaults.sh"
