# rclone/obsidian_vaults/driveBiSync_obsidian-vaults.cron
@reboot ~/dotfiles/rclone/obsidian_vaults/driveBiSync_obsidian-vaults.sh
*/15 * * * * ~/dotfiles/rclone/obsidian_vaults/driveBiSync_obsidian-vaults.sh
