# rclone/driveBiSync.cron
@reboot ~/dotfiles/rclone/driveBiSync.sh
*/10 * * * * ~/dotfiles/rclone/driveBiSync.sh
